# root/etc/s6-overlay/s6-rc.d/init-jellyfin-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# create directories
mkdir -p \
    /config/{log,data/plugins/configurations,data/transcodes,cache} \
    /data \
    /transcode

# permissions
lsiown abc:abc \
    /config \
    /config/* \
    /config/data/plugins \
    /config/data/plugins/configurations \
    /config/data/transcodes \
    /transcode
